# udev/.config/udev-rules/98-hdmi-detect.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="/home/johnny/.local/bin/hdmiconnected"
