# autostart/.config/autostart/unclutter.desktop
[Desktop Entry]
Name=Unclutter
Exec=unclutter -idle 0 
Type=Application
X-GNOME-Autostart-enabled=true

